# services/macos/com.williamveith.fbs-interlock-gateway.plist.in
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>@MACOS_LAUNCHD_LABEL@</string>

    <key>ProgramArguments</key>
    <array>
        <string>@MACOS_INSTALL_DIR@/start.sh</string>
    </array>

    <key>WorkingDirectory</key>
    <string>@MACOS_INSTALL_DIR@</string>

    <key>UserName</key>
    <string>@MACOS_SERVICE_USER@</string>

    <key>GroupName</key>
    <string>@MACOS_SERVICE_GROUP@</string>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <true/>

    <key>ThrottleInterval</key>
    <integer>2</integer>

    <key>ProcessType</key>
    <string>Background</string>

    <key>StandardOutPath</key>
    <string>@MACOS_LOG_DIR@/gateway.log</string>

    <key>StandardErrorPath</key>
    <string>@MACOS_LOG_DIR@/gateway-error.log</string>
</dict>
</plist>
