# services/linux/app.service.in
[Unit]
Description=FBS Interlock Gateway
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=@SERVICE_USER@
Group=@SERVICE_GROUP@
WorkingDirectory=@INSTALL_DIR@
ExecStart=@INSTALL_DIR@/@APP@ -config @CONFIG_PATH@
Restart=always
RestartSec=2
StartLimitIntervalSec=60
StartLimitBurst=10
NoNewPrivileges=true

[Install]
WantedBy=multi-user.target